# openclaw-telegram-wsl-setup/tools/openclaw-gateway-recovery/openclaw-gateway-recovery.service
[Unit]
Description=OpenClaw Gateway Recovery
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/openclaw-gateway-recovery.env
ExecStart=%h/.local/bin/openclaw-gateway-recovery
Nice=10
