# openclaw-telegram-wsl-setup/tools/agent-room/openclaw-agent-room-retention-prune.timer
[Unit]
Description=Run OpenClaw Agent Room retention prune hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
