# modules/openclaw-market-immersion/systemd/openclaw-market-immersion-close.service
[Unit]
Description=OpenClaw Market Immersion close run
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
Restart=on-failure
RestartSec=2min
RestartPreventExitStatus=3
WorkingDirectory=%h/.openclaw/workspace/market-immersion-module
Environment=MARKET_IMMERSION_DIR=%h/.openclaw/workspace/market-immersion-module
EnvironmentFile=%h/.openclaw/secrets/mx.env
Environment=MARKET_IMMERSION_PHASE=close
ExecStart=%h/.openclaw/workspace/market-immersion-module/scripts/run_market_immersion.sh close --resume
