# modules/openclaw-market-immersion/systemd/openclaw-market-feed-snapshot.timer
[Unit]
Description=Schedule OpenClaw Market Feed Snapshot retention

[Timer]
OnCalendar=*:00/30
Persistent=true
RandomizedDelaySec=60
Unit=openclaw-market-feed-snapshot.service

[Install]
WantedBy=timers.target
