# modules/openclaw-market-immersion/systemd/openclaw-market-immersion-close.timer
[Unit]
Description=Schedule OpenClaw Market Immersion close collection

[Timer]
OnCalendar=*-*-* 15:20:00
Persistent=true
Unit=openclaw-market-immersion-close.service

[Install]
WantedBy=timers.target
