# modules/openclaw-market-immersion/systemd/openclaw-market-feed-snapshot.service
[Unit]
Description=OpenClaw Market Feed Snapshot retention run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.openclaw/workspace/market-immersion-module
Environment=MARKET_IMMERSION_DIR=%h/.openclaw/workspace/market-immersion-module
EnvironmentFile=%h/.openclaw/secrets/mx.env
ExecStart=%h/.openclaw/workspace/.venv-mxskills/bin/python %h/.openclaw/workspace/market-immersion-module/scripts/market_feed_snapshot.py --lookback-minutes 90
TimeoutStartSec=600
Restart=no
