# modules/openclaw-market-immersion/systemd/openclaw-source-interface-verification.timer
[Unit]
Description=Run OpenClaw Market Source Interface Verification twice monthly

[Timer]
OnBootSec=10min
OnCalendar=*-*-01,16 07:05:00
AccuracySec=1h
Persistent=true
Unit=openclaw-source-interface-verification.service

[Install]
WantedBy=timers.target
