# modules/openclaw-market-immersion/systemd/openclaw-market-immersion-night.timer
[Unit]
Description=Schedule OpenClaw Market Immersion night collection

[Timer]
OnCalendar=*-*-* 22:10:00
Persistent=true
Unit=openclaw-market-immersion-night.service

[Install]
WantedBy=timers.target

