# openclaw-telegram-wsl-setup/tools/openclaw-netwatch/openclaw-netwatch.service
[Unit]
Description=OpenClaw network recovery watchdog
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/openclaw-netwatch.env
ExecStart=%h/.local/bin/openclaw-netwatch
